# srtla.service
[Unit]
Description=srtla
 
[Service]
ExecStart=/bin/bash /home/ubuntu/srtla.sh

[Install]
WantedBy=multi-user.target
